# installer/bridge/com.rushautoworks.racestudio3.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- aim-bridge root LaunchDaemon, registered via SMAppService from RaceStudio 3.app.
     Runs as root (exempt from the macOS Local Network gate) and relays loopback <-> the AiM dash.
     BundleProgram is relative to the app bundle; the daemon hard-pins 10.0.0.1:2000/36002 when
     running as root (it ignores the environment), so no arguments carry a destination. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rushautoworks.racestudio3.bridge</string>
    <key>BundleProgram</key>
    <string>Contents/MacOS/aim-bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/aim-bridge.log</string>
    <key>StandardOutPath</key>
    <string>/Library/Logs/aim-bridge.log</string>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.rushautoworks.racestudio3</string>
    </array>
</dict>
</plist>
